# savemoney/supervisord.conf
; supervisor config file
[inet_http_server]
port=127.0.0.1:9001
username=user
password=123

[supervisorctl]
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0777                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf




[program:uwsgi-savemoney]
environment=PATH="/home/savemoneyenv2/bin",VIRTUAL_ENV="/home/savemoneyenv2/savemoney/",PYTHONPATH="/home/savemoneyenv2/lib/python:/home/savemoneyenv2/lib/python3.6/site-packages"
# 确保uwsgi命令的路径是对的
directory=/home/savemoneyenv2/savemoney
command=uwsgi --ini uwsgi.ini


numprocs=1          ;进程数
autostart=true      ;当supervisor启动时,程序将会自动启动 
autorestart=true    ;自动重启
stopsignal=INT

# 注意确保路径存在
stdout_logfile=/home/savemoneyenv2/savemoney/site_out.log
stderr_logfile=/home/savemoneyenv2/savemoney/site_err.log
stopsignal=QUIT
killasgroup=true


[program:savemoney_celery]
environment=PATH="/home/savemoneyenv2/bin",VIRTUAL_ENV="/home/savemoneyenv2/savemoney/",PYTHONPATH="/home/savemoneyenv2/lib/python:/home/savemoneyenv2/lib/python3.6/site-packages"
command=celery -A savemoney worker --loglevel info
directory=/home/savemoneyenv2/savemoney/


user=root
; Supervisor will start as many instances of this program as named by numprocs
numprocs=1
stdout_logfile=/home/savemoneyenv2/savemoney/celery_worker.log 
stderr_logfile=/home/savemoneyenv2/savemoney/celery_worker_err.log 
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true
priority=998

[program:celery.beat] 
environment=PATH="/home/savemoneyenv2/bin",VIRTUAL_ENV="/home/savemoneyenv2/savemoney/",PYTHONPATH="/home/savemoneyenv2/lib/python:/home/savemoneyenv2/lib/python3.6/site-packages"
;指定运行目录 
directory=/home/savemoneyenv2/savemoney/
;运行目录下执行命令
command=celery -A savemoney beat --loglevel info
;运行目录下执行命令

user=ivan
;启动设置 
numprocs=1          ;进程数
autostart=true      ;当supervisor启动时,程序将会自动启动 
autorestart=true    ;自动重启
 
;停止信号
stopsignal=INT


;输出日志 
stdout_logfile=/home/savemoneyenv2/savemoney/celery_beat.log 
stdout_logfile_maxbytes=10MB  ;默认最大50M 
stdout_logfile_backups=10     ;日志文件备份数，默认为10 
 
;错误日志 
redirect_stderr=false         ;为true表示禁止监听错误 
stderr_logfile=/home/savemoneyenv2/savemoney/celery_beat_err.log 
stderr_logfile_maxbytes=10MB 
stderr_logfile_backups=10



